# odroidu3-fanspeed-autoctrl.service
[Unit]
Description=ODROID U3 fanspeed control

[Service]
ExecStart=/usr/local/bin/odroidu3-fanspeed-autoctrl

[Install]
WantedBy=multi-user.target
